# server/deploy/usb-control.service
[Unit]
Description=USB Control Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/usb-control/bin/usb-control --config /etc/usb-control/usb-control.toml
Restart=on-failure
RestartSec=3
WorkingDirectory=/opt/usb-control
Environment=RUST_LOG=info
LimitNOFILE=65536
RuntimeDirectory=usb-control
RuntimeDirectoryMode=0700

[Install]
WantedBy=multi-user.target
